# AnIA/systemd/ania-api.service
[Unit]
Description=AnIA FastAPI backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/proyectos/DAW/AnIA
EnvironmentFile=/home/%i/proyectos/DAW/AnIA/.env
ExecStart=/home/%i/proyectos/DAW/AnIA/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
